# quickshell/quickshell-reset-app-usage.service
[Unit]
Description=Reset QuickShell launcher app usage on session shutdown
Documentation=file:%h/.config/quickshell/app_usage.json

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/true
ExecStop=%h/.config/quickshell/scripts/reset-app-usage.sh

[Install]
WantedBy=default.target
